# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null

[program:apache]
command=bash -c "php /app/artisan config:clear && php /app/artisan cache:clear && a2dismod mpm_event mpm_worker 2>/dev/null; a2enmod mpm_prefork rewrite 2>/dev/null; apache2-foreground"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ml_server]
command=python3 /app/ml_server.py
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0